# ssserver.ini
[program:ssserver]
process_name=%(program_name)s_%(process_num)02d
command=ssserver -c /etc/shadowsocks.json
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/root/ssserver.log
